# fleet/systemd/agv-fleet-manager.service
# AGV Fleet Manager — VDA 5050 Master Control
# Install: sudo cp agv-fleet-manager.service /etc/systemd/system/
# Enable:  sudo systemctl enable agv-fleet-manager
# Start:   sudo systemctl start agv-fleet-manager

[Unit]
Description=AGV Fleet Manager
After=network.target agv-mosquitto.service

[Service]
Type=simple
User=orza
WorkingDirectory=/home/orza/ros2_ws/fleet/agv_fleet_manager
ExecStart=/usr/bin/node dist/index.js
Environment="VDA_MQTT_BROKER=mqtt://localhost:1883"
# 8091 is taken by agv_image_server on the same Jetson (specs/project.yaml).
Environment="FLEET_PORT=8092"
# Shared secret for the REST/WebSocket control surface (fleet e-stop,
# navigate, traffic zones). Unset = unauthenticated; isolated LAN only.
#Environment="FLEET_API_TOKEN=change-me"
# Broker credentials, once mosquitto auth is enabled (see mosquitto.conf):
#Environment="VDA_MQTT_USERNAME=fleet-manager"
#Environment="VDA_MQTT_PASSWORD=change-me"
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
